# luci-theme-argon-gray/root/etc/uci-defaults/30_luci-theme-argon-gray
#!/bin/sh

/sbin/uci set luci.themes.ArgonGray='/luci-static/argon-gray'
/sbin/uci set luci.main.mediaurlbase='/luci-static/argon-gray'
/sbin/uci commit luci

exit 0
